# deploy/launchd/com.nba.closing_odds.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nba.closing_odds</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/shawnclaw/autobot/investing/sports/autobots_NBA/fetch_closing_odds.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/shawnclaw/autobot/investing/sports/autobots_NBA</string>

  <!-- NBA 開賽集中在台灣時間 01/03/07-10 點，各時段於開賽前 30-60 分鐘抓快照 -->
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key><integer>0</integer>
      <key>Minute</key><integer>30</integer>
    </dict>
    <dict>
      <key>Hour</key><integer>2</integer>
      <key>Minute</key><integer>30</integer>
    </dict>
    <dict>
      <key>Hour</key><integer>6</integer>
      <key>Minute</key><integer>30</integer>
    </dict>
    <dict>
      <key>Hour</key><integer>7</integer>
      <key>Minute</key><integer>30</integer>
    </dict>
    <dict>
      <key>Hour</key><integer>8</integer>
      <key>Minute</key><integer>30</integer>
    </dict>
    <dict>
      <key>Hour</key><integer>10</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/shawnclaw/autobot/investing/sports/autobots_NBA/logs/closing_odds.launchd.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/shawnclaw/autobot/investing/sports/autobots_NBA/logs/closing_odds.launchd.error.log</string>

  <key>SoftResourceLimits</key>
  <dict>
    <key>NumberOfFiles</key>
    <integer>4096</integer>
  </dict>

  <key>HardResourceLimits</key>
  <dict>
    <key>NumberOfFiles</key>
    <integer>8192</integer>
  </dict>
</dict>
</plist>
